# scripts/dpkg-deb/postrm
USER=lorawan

case "$1" in
    remove|purge)
    systemctl disable lorawan-server
    rm -rf /var/log/lorawan-server
    userdel -r $USER
    ;;

    *)
    ;;
esac
