# scripts/dpkg-deb/postrm
USER=lorawan

case "$1" in
    remove|purge)
    userdel -r $USER
    systemctl disable lorawan-server
    ;;

    *)
    ;;
esac
